# supervisord.conf
[supervisord]
nodaemon=true

[program:voicevox]
command=/usr/local/bin/voicevox_engine
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:app]
command=python3 /app/main.py
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
